# contrib/coreos/systemd/fleet.service
[Unit]
Description=fleet

[Service]
ExecStart=/usr/bin/fleet

[Install]
WantedBy=multi-user.target