# controller/systemd/deis-controller.service
[Unit]
Description=deis-controller

[Service]
TimeoutStartSec=20m
ExecStartPre=/usr/bin/docker pull deis/controller:latest
ExecStart=/usr/bin/docker run --name deis-controller -p 8000:8000 -e PUBLISH=8000 deis/controller
ExecStop=/usr/bin/docker rm -f deis-controller

[Install]
WantedBy=multi-user.target
