# controller/systemd/deis-controller.service
[Unit]
Description=deis-controller
After=docker.service
Requires=docker.service
Requires=docker-tcp.socket

[Service]
ExecStart=/bin/bash -c '/usr/bin/docker start -a deis-controller || /usr/bin/docker run --name deis-controller -P -e HOST=10.1.42.1 -e ETCD=10.1.42.1:4001 deis/controller'
ExecStop=/usr/bin/docker stop deis-controller

[Install]
WantedBy=multi-user.target