# logger/systemd/deis-logger.service
[Unit]
Description=deis-logger

[Service]
TimeoutStartSec=20m
ExecStartPre=/usr/bin/docker pull deis/logger:latest
ExecStart=/usr/bin/docker run --name deis-logger -p 514:514 -e PUBLISH=514 deis/logger
ExecStop=/usr/bin/docker rm -f deis-logger

[Install]
WantedBy=multi-user.target
