# builder/systemd/deis-builder.service
[Unit]
Description=deis-builder

[Service]
TimeoutStartSec=20m
ExecStartPre=/usr/bin/docker pull deis/builder:latest
ExecStart=/usr/bin/docker run --name deis-builder -p 2222:22 -e PUBLISH=22 -e PORT=2222 --privileged deis/builder
ExecStop=/usr/bin/docker rm -f deis-builder

[Install]
WantedBy=multi-user.target
